# services/orggcalsync.service
[Unit]
Description=Call orgsync
Wants=orgcalsync.timer

[Service]
ExecStart=/home/japhir/bin/orggcalsync

[Install]
WantedBy=timers.target